# deploy/systemd/designo-backend.service
[Unit]
Description=Designo Backend (Motion Website Generator API)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/app/designo/backend
ExecStart=/opt/app/designo/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8620
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
